# trading_llm_bot_simulation.service
[Unit]
Description=Trading LLM Bot Simulation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
Group=%GROUP%
WorkingDirectory=%PROJECT_ROOT%
Environment="PATH=%VENV_BIN%:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=%PROJECT_ROOT%"
# Load environment variables from .env file if it exists
EnvironmentFile=-%PROJECT_ROOT%/.env

# Gunicorn command for simulation service
ExecStart=%VENV_BIN%/gunicorn -c simulate_gunicorn.conf.py simulate_wsgi:application

# Graceful reload
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy - restart on failure to ensure simulation keeps running
# Longer delay to allow network/DNS to stabilize
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading_llm_bot_sim

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
